# jc-ldap-proxy.container
# /etc/containers/systemd/jc-ldap-proxy.container
# Quadlet unit — after editing: systemctl daemon-reload && systemctl start jc-ldap-proxy
[Unit]
Description=JumpCloud LDAP caching proxy for QNAP
After=network-online.target
Wants=network-online.target

[Container]
# Replace with an immutable release tag or manifest digest before installing.
Image=ghcr.io/theaustrian75/jumpcloud-ldap-proxy@sha256:REPLACE_WITH_MANIFEST_DIGEST
ContainerName=jc-ldap-proxy
# QNAP httpdusr:everyone defaults. Quadlet does not interpolate .env values;
# change User, Group, and both Tmpfs uid/gid options together if needed.
User=99
Group=100
Environment=JC_ORG_ID=CHANGE_ME
# UID of the JumpCloud account configured as QTS's LDAP bind user
Environment=JC_CACHE_READER_UID=CHANGE_ME
Environment=SLAPD_LOGLEVEL=stats
Environment=TZ=America/New_York
# Uncomment if a remote client must reach this proxy (IP or ip%netmask)
#Environment=ALLOWED_CLIENT_IP=192.168.1.50
PublishPort=127.0.0.1:636:636
PublishPort=127.0.0.1:389:389
# Server cert/key for the LDAPS listener, issued by your lab CA
# (SAN must match the name/IP the QNAP will use)
Volume=/srv/jc-ldap-proxy/certs:/certs:ro,Z
# Runtime files and cached NTLM hashes live on tmpfs and never touch disk.
# Cost: cache starts cold after each restart (first lookups pay one WAN
# round-trip again).
Tmpfs=/run/openldap:rw,size=16m,uid=99,gid=100,mode=0700
Tmpfs=/var/lib/ldap/pcache:rw,size=256m,uid=99,gid=100,mode=0700
ReadOnly=true
NoNewPrivileges=true
DropCapability=all
AddCapability=CAP_NET_BIND_SERVICE
PidsLimit=128
Memory=512m
StopTimeout=15
HealthOnFailure=restart

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
